# examples/http-activation/http-activation-b.socket
[Unit]
Description=HTTP Activation Socket B

[Socket]
ListenStream=10081
Service=http-activation.service
FileDescriptorName=b
