# examples/http-activation/http-activation-a.socket
[Unit]
Description=HTTP Activation Socket A

[Socket]
ListenStream=10080
Service=http-activation.service
FileDescriptorName=a
