# examples/notify/notify.service
[Unit]
Description=Notify Service

[Service]
ExecStart=/home/vic/bin/notify
Type=notify
